# netwatcher.service
[Unit]
Description=NetWatcher - Network Packet Monitoring System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netwatcher
EnvironmentFile=/opt/netwatcher/.env
Environment=LANG=en_US.UTF-8 LC_ALL=en_US.UTF-8 PYTHONIOENCODING=utf-8
ExecStart=/opt/netwatcher/.venv/bin/python -m netwatcher
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=no
ProtectSystem=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
